# raspberry/config/systemd/neopro-score-bridge.service
# Systemd service pour le Score Bridge Neopro
# Écoute les mises à jour de score via Socket.IO et écrit dans des fichiers
# pour l'incrustation FFmpeg

[Unit]
Description=Neopro Score Bridge (Socket.IO to Files)
Documentation=https://github.com/neopro
After=network.target neopro-app.service
Wants=neopro-app.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/neopro/services
Environment=NODE_ENV=production
Environment=SOCKET_URL=http://localhost:3000
Environment=SCORE_DIR=/tmp

# Commande de lancement
ExecStart=/usr/bin/node /home/pi/neopro/services/score-bridge.js

# Redémarrage automatique en cas de crash
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neopro-score-bridge

# Limites de sécurité
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
